# ops/systemd/captiva-mail-gateway.service
[Unit]
Description=CAPTIVA Mail Gateway
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/CAPTIVA
EnvironmentFile=/var/www/CAPTIVA/ops/mail-gateway/.env
ExecStart=/usr/bin/node /var/www/CAPTIVA/ops/mail-gateway/server.mjs
Restart=always
RestartSec=5
User=www-data
Group=www-data
StandardOutput=append:/var/log/captiva/mail-gateway.log
StandardError=append:/var/log/captiva/mail-gateway.err.log

[Install]
WantedBy=multi-user.target
