# ops/systemd/myoptiwealth-saas-healthcheck.service
[Unit]
Description=MyOptiwealth SaaS healthcheck
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/var/www/myoptiwealth
EnvironmentFile=-/etc/myoptiwealth-saas/healthcheck.env
ExecStart=/var/www/myoptiwealth/ops/saas-healthcheck.sh
