# ops/systemd/qrt-ops-worker.service
[Unit]
Description=MYOPTIWEALTH QRT Ops Worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/myoptiwealth
EnvironmentFile=/var/www/myoptiwealth/.env
ExecStart=/bin/bash -lc 'while true; do /usr/bin/npm run ops:qrt:once; sleep 30; done'
Restart=always
RestartSec=5
User=www-data
Group=www-data
StandardOutput=append:/var/log/myoptiwealth/qrt-ops-worker.log
StandardError=append:/var/log/myoptiwealth/qrt-ops-worker.err.log

[Install]
WantedBy=multi-user.target
